# .config/systemd/user/offlineimap.service
[Unit]
Description=OfflineIMAP Service

[Service]
Type=oneshot
Environment=NOTMUCH_CONFIG=%h/.config/notmuch/notmuchrc
ExecStart=/usr/bin/offlineimap -c %h/.local/private/offlineimap/offlineimaprc
ExecStartPost=/usr/bin/notmuch new

[Install]
WantedBy=default.target
